# systemd-user/.config/systemd/user/clamav-hourly.timer
[Unit]
Description=Run ClamAV scan every hour (user)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
