# systemd-user/.config/systemd/user/clamav-hourly.service
[Unit]
Description=User ClamAV Hourly Scan

[Service]
Type=oneshot
ExecStart=%h/.local/bin/clamav-hourly.sh
