# systemd-user/.config/systemd/user/hyprsunset-on.service
[Unit]
Description=Start hyprsunset night filter

[Service]
Type=simple
ExecStart=/usr/bin/hyprsunset -t 3500
Restart=on-failure
