# systemd-user/.config/systemd/user/hyprsunset-off.service
[Unit]
Description=Stop hyprsunset night filter

[Service]
Type=oneshot
ExecStart=/usr/bin/killall hyprsunset
